# deploy/wloc/pdg-wloc.service
[Unit]
Description=PrivDNS Gateway scoped shared MITM sidecar
After=network-online.target
Wants=network-online.target
ConditionPathExists=/opt/pdg-bot/wloc_mitm.py
ConditionPathExists=/opt/pdg-bot/adblock_mitm.py

[Service]
Type=simple
User=pdg-wloc
Group=pdg-wloc
Environment=HOME=/var/lib/pdg-wloc
WorkingDirectory=/var/lib/pdg-wloc
ExecStart=/usr/bin/mitmdump --mode regular --listen-host 127.0.0.1 --listen-port 9080 --set confdir=/var/lib/pdg-wloc/mitmproxy --set connection_strategy=lazy --set upstream_cert=false --set block_global=false --scripts /opt/pdg-bot/wloc_mitm.py --scripts /opt/pdg-bot/adblock_mitm.py
Restart=on-failure
RestartSec=3
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/pdg-wloc
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
MemoryMax=256M

[Install]
WantedBy=multi-user.target
