# deploy/bot/pdg-rules-update.service
[Unit]
Description=PrivDNS Gateway 定时刷新 geosite + 分流规则集 + 去广告规则
After=network-online.target mihomo.service mosdns.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/pdg-bot/scheduled-update.sh
